# root/etc/uci-defaults/99-haproxy
#!/bin/sh

uci -q get haproxy.globals > /dev/null 2>&1 || {
    uci batch <<-EOF
        set haproxy.globals=globals
        set haproxy.globals.maxconn=4096
        set haproxy.globals.log='stdout local0'
        set haproxy.globals.stats_socket=/var/run/haproxy/haproxy.sock
        set haproxy.globals.inspect_delay=5s
        commit haproxy
EOF
}

# Enable boot start for both the config generator (runs first, builds the
# config + creates the stats-socket dir) and haproxy itself. Neither is
# enabled by default on install, so without this haproxy never comes up at boot.
/etc/init.d/haproxy-gen enable 2>/dev/null
/etc/init.d/haproxy     enable 2>/dev/null

# Generate the config now, then (re)start haproxy so it is live immediately
# after install rather than only after the next reboot.
/usr/sbin/haproxy-gen 2>/dev/null || true
/etc/init.d/haproxy restart 2>/dev/null || /etc/init.d/haproxy start 2>/dev/null || true

exit 0
